# docker/kubeflow/s6/services.d/opencode/run
#!/command/with-contenv sh

# Run the combined OpenCode API + UI server
exec /usr/local/bin/start-server.sh
